# djangoq.socket
[Unit]
Description=djangoq socket

[Socket]
ListenStream=/run/djangoq.sock

[Install]
WantedBy=sockets.target